# deploy/systemd/lemma-live-publisher.service
[Unit]
Description=Publish Lemma live task JSON
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/lemma
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/.local/bin/uv run lemma dashboard publish --output-dir /var/www/lemma-live
