# deploy/systemd/lemma-public-dashboard.service
[Unit]
Description=Publish Lemma public dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/lemma
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=LEMMA_SRC=/opt/lemma
Environment=EXPORTER_SRC=/opt/lemma
Environment=SITE_SRC=/opt/lemmasub.net
Environment=SUMMARY_JSONL=/var/lib/lemma/public-summary.jsonl
Environment=NETUID=467
Environment=SUBTENSOR_NETWORK=test
Environment=SUBTENSOR_CHAIN_ENDPOINT=wss://test.finney.opentensor.ai:443
ExecStart=/opt/lemma/deploy/scripts/lemma-refresh-public-dashboard
